# cpu_scheduler/perf-prof.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/perf-prof.sh

[Install]
WantedBy=default.target